# rootfs/etc/services.d/bind/run
#!/usr/bin/with-contenv sh

exec named -f -c /config/etc/named.conf
